# debian/extra/rules/80-networking.rules
SUBSYSTEM=="net", ACTION=="add|remove", RUN+="net.agent"

SUBSYSTEM=="atm", ACTION=="add|remove", RUN+="dsl-modem.agent"

