# debian/extra/rules/80-networking.rules
SUBSYSTEM=="net", ACTION=="add|remove", RUN+="net.agent"
