# debian/system/halt.target
[Unit]
Description=Halt
Requires=halt.service
