# debian/system/multi-user.target
[Unit]
Description=Multi-User Target
Requires=basic.target
After=basic.target
Names=runlevel2.target
