# debian/system/nss-lookup.target
[Unit]
Description=Name Lookups
Wants=network.target
After=network.target
