# debian/system/poweroff.target
[Unit]
Description=Power Off
Names=runlevel0.target
Requires=poweroff.service
