# debian/system/reboot.target
[Unit]
Description=Reboot
Names=runlevel6.target
Requires=reboot.service
