# debian/system/rescue.target
[Unit]
Description=Rescue Target
Requires=basic.target
After=basic.target
Names=runlevel1.target
