# debian/system/runlevel4.target
[Unit]
Description=Compatibility for Runlevel 4
Requires=multi-user.target
After=multi-user.target
