# debian/system/runlevel5.target
[Unit]
Description=Compatibility for Runlevel 5
Requires=multi-user.target
After=multi-user.target
