# debian/system/shutdown.target
[Unit]
Names=shutdown.target
Description=Shutdown
