# debian/system/sysinit.service
[Unit]
Description=System Initialization

[Service]
ExecStart=/etc/init.d/rcS
Type=finish
TimeoutSec=0
StandardInput=tty
ValidNoProcess=yes
