# debian/system/syslog.target
[Unit]
Description=System Logger
Requires=rsyslog.service
After=rsyslog.service
