# debian/extra/units/ifup@.service
[Unit]
Description=ifup for %I
After=local-fs.target network-pre.target systemd-tmpfiles-setup.service
Before=network.target
BindsTo=sys-subsystem-net-devices-%i.device
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/lib/systemd/systemd-ifupdown-helper start %I
ExecStop=/lib/systemd/systemd-ifupdown-helper stop %I
RemainAfterExit=true
