# debian/extra/units/networking.service.d/systemd.conf
[Unit]
After=network-pre.target apparmor.service
